# scripts/jarvit-vm-update.timer
[Unit]
Description=JARVIT VM Auto-Update Timer
Documentation=https://github.com/jarvit242351QE/jarvit-vm

[Timer]
# Run every 30 minutes, with a random delay up to 5 minutes
# to avoid all VMs hitting GitHub API at the exact same time
OnBootSec=2min
OnUnitActiveSec=30min
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
